# src/etc/init.d/one-context-force##apk.one
#!/sbin/openrc-run

name="OpenNebula forced run"

depend() {
	use one-context one-context-local
	keyword -stop -shutdown
}

start() {
	ebegin "Starting one-context-force"
	export TIMEOUT=0
	/usr/sbin/one-context-run force
	eend $?
}
