# src/etc/init.d/one-context-local##apk
#!/sbin/openrc-run

name="OpenNebula pre-networking contextualization"

depend() {
	need localmount udev udev-trigger udev-settle
	before net iptables keepalived
	use logger
	keyword -stop -shutdown
}

start() {
	ebegin "Starting one-context-local"
	/usr/sbin/one-contextd local
	eend $?
}
