# src/etc/init.d/one-context-reconfigure-delayed##apk
#!/sbin/openrc-run

name="OpenNebula delayed reconfiguration"

depend() {
	need one-context one-context-local
	keyword -stop -shutdown
}

start() {
	ebegin "Starting one-context-reconfigure-delayed"
	/usr/sbin/one-context-reconfigure
	eend $?
}
